# cron.daily/iplayer
#!/usr/bin/env bash
set -o errexit
sudo -u pi /home/pi/bin/grab
