# cron.daily/minidlna
#!/bin/sh
service minidlna force-reload
